# .github/workflows/OnlineLLM.service
[Unit]
Description=OnlineLLM
After=network.target

[Service]
ExecStart=/usr/local/bin/OnlineLLM
Restart=always
User=ubuntu
Environment="OPENAI_API_KEY=%OPENAI_API_KEY%"
WorkingDirectory=/home/ubuntu
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target